# command/com.fabrizio.fattura-watchers.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
   "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <!-- Define a unique label -->
    <key>Label</key>
    <string>com.fabrizio.fattura.watchers</string>

    <!-- Set the working directory to /Users/fabrizio -->
    <key>WorkingDirectory</key>
    <string>/Users/fabrizio</string>

    <!-- Command to execute: the shell script to run -->
    <key>ProgramArguments</key>
    <array>
      <string>/Users/fabrizio/DSLegno/command/fattura-watchers.sh</string>
    </array>

    <!-- Start the job as soon as you log in -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart the process if it crashes -->
    <key>KeepAlive</key>
    <true/>

    <!-- Log output to a single file for debugging purposes -->
    <key>StandardOutPath</key>
    <string>/Users/fabrizio/log.txt</string>
    <key>StandardErrorPath</key>
    <string>/Users/fabrizio/log.txt</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
      <key>HOME</key>
      <string>/Users/fabrizio</string>
      <!-- Ensure PATH includes the directory for watchexec -->
      <key>PATH</key>
      <string>/opt/homebrew/bin</string>
    </dict>
  </dict>
</plist>
